# GentooGuide/files/etc/profile.d/marker.sh
#!/usr/bin/env bash
# ############################################################################
# # PATH: /etc/profile.d/                         AUTHOR: Hoefkens.j@gmail.com
# # FILE: marker.sh
# ############################################################################
#
[[ -n $LOADED ]] && LOADED=() 
LOADED+=("/etc/profile.d")

